# roles/supervisord/templates/supervisord.conf.j2
[supervisord]
nodeamon=false

[program:test]
command=php /var/www/{{name_server}}/artisan start
